# templates/game-sync.timer
[Unit]
Description=Run Game sync every 1 minute

[Timer]
OnActiveSec=1min
OnUnitInactiveSec=1min

[Install]
WantedBy=timers.target
