# templates/game-sync.timer
[Unit]
Description=Run Game sync every 5 minutes

[Timer]
OnActiveSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
